# setup/brian.service
[Unit]
Description=Brian AI Agent
After=network.target docker.service

[Service]
Type=simple
User=brian
WorkingDirectory=/home/brian/app
EnvironmentFile=/etc/brian/env
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target
